# infra/systemd/apisix.service
[Unit]
Description=Apache APISIX Gateway (Docker container)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Restart=always
RestartSec=10s
TimeoutStartSec=120
TimeoutStopSec=30

ExecStartPre=-/usr/bin/docker rm -f messenger-apisix
ExecStartPre=/usr/bin/docker network create messenger_net 2>/dev/null || true

ExecStart=/usr/bin/docker run --rm \
    --name messenger-apisix \
    --network messenger_net \
    -p 80:9080 \
    -p 443:9443 \
    -v /opt/messenger/repo/infra/apisix/config.yaml:/usr/local/apisix/conf/config.yaml:ro \
    -v /opt/messenger/repo/infra/apisix/apisix.yaml:/usr/local/apisix/conf/apisix.yaml:ro \
    -v /etc/letsencrypt:/etc/letsencrypt:ro \
    -v /opt/messenger/repo/infra/nginx/downloads:/usr/share/nginx/html/downloads:ro \
    -v /opt/messenger/certbot-webroot:/var/www/certbot:ro \
    apache/apisix:3.15.0-debian

ExecStop=/usr/bin/docker stop messenger-apisix

StandardOutput=journal
StandardError=journal
SyslogIdentifier=apisix

[Install]
WantedBy=multi-user.target
